# deploy/systemd/betbot-ml.service
[Unit]
Description=betbot stub service

[Service]
ExecStart=/usr/bin/true
